# configs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# ---------------------------------------------------------------------------
# OpenClaw upstream gateway
# ---------------------------------------------------------------------------
# Command mirrors the upstream image's default CMD (docker inspect
# ghcr.io/openclaw/openclaw:2026.4.14): `node openclaw.mjs gateway
# --allow-unconfigured`. The --allow-unconfigured flag lets the gateway
# start before post_up hooks apply config; dropping it would deadlock
# the first deploy (container refuses to start, hooks never run).
# Skips upstream's docker-entrypoint.sh shim (a `set -- node "$@"`
# wrapper) since supervisord is already PID 1.
[program:openclaw]
command=node /app/openclaw.mjs gateway --allow-unconfigured
directory=/app
user=node
# HOME/USER must be set explicitly: supervisord's user= directive only changes
# uid — HOME/USER env vars are inherited from supervisord (root). openclaw reads
# $HOME/.openclaw for its state; without this it tries /root/.openclaw and
# fails EACCES. See volume mounts in the module compose fragment — they target
# /home/node/.openclaw (v2026.4.14-2).
environment=HOME="/home/node",USER="node"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# ---------------------------------------------------------------------------
# Primare shim (FastAPI /info /health /metrics on port 8000)
# ---------------------------------------------------------------------------
[program:shim]
command=/opt/shim/.venv/bin/uvicorn primare_module_openclaw.main:app --host 0.0.0.0 --port 8000
directory=/opt/shim
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
