# configs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# ---------------------------------------------------------------------------
# OpenClaw upstream gateway
# ---------------------------------------------------------------------------
# TODO: verify entrypoint on first build. Inspect the upstream image with:
#   docker inspect ghcr.io/openclaw/openclaw:2026.4.14 --format '{{json .Config.Entrypoint}} {{json .Config.Cmd}}'
# Update the command= line below to match. Current best-guess based on the
# upstream project CLI is `openclaw gateway start`; adjust if needed.
[program:openclaw]
command=/usr/local/bin/openclaw gateway start
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# ---------------------------------------------------------------------------
# Primare shim (FastAPI /info /health /metrics on port 8000)
# ---------------------------------------------------------------------------
[program:shim]
command=/opt/shim/.venv/bin/uvicorn primare_module_openclaw.main:app --host 0.0.0.0 --port 8000
directory=/opt/shim
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
